# userspace/spider-rs/units/basic.target
[Unit]
Name=Basic Target
Description=Mirage basic system target

[Target]
