# userspace/spider-rs/units/default.target
[Unit]
Name=Default Target
Description=Mirage default boot target
Wants=basic.target m1-terminal.service
After=basic.target

[Target]
