# userspace/spider-rs/units/emergency.target
[Unit]
Description=Emergency Spider target
